# Windows/ScheduledTask_Interpreter.yaml
title: Scheduled Task Privilege Escalation Through an Interpreter
id: b3f7e8d2-4c1a-4e9b-8d6f-1a5e7c9b2d4a
status: stable
description: Detects suspicious scheduled task creation for privilege escalation
references:
    - https://attack.mitre.org/techniques/T1053/005/
    - https://thehackernews.com/2025/04/experts-uncover-four-new-privilege.html
author: zendanny
date: 2026/03/09
tags:
    - attack.privilege-scalation
    - attack.persistence
    - attack.t1053.005
logsource:
    product: windows
    service: security
detection:
    selection_task_creation:
        EventID: 4698  # Scheduled task created
        TaskName|regex: '^(?!.*\\Microsoft\\Windows\\).*'  # Exclude Microsoft tasks
    selection_high_privilege_indicators:
        SubjectUserName|endswith: '$'  # Computer account creating task
        TaskContent|contains:
            - '<RunLevel>HighestAvailable</RunLevel>'
            - '<UserId>S-1-5-18</UserId>'  # SYSTEM
            - '<LogonType>S4U</LogonType>'  # Service for user
    selection_suspicious_actions:
        TaskContent|contains:
            - 'cmd.exe /c'
            - 'cscript.exe'
            - 'wscript.exe'
            - 'powershell.exe -'
    selection_suspicious_paths:
        TaskContent|contains:
            - '\AppData\'
            - '\Users\Public\'
            - '\Temp\'
            - '\ProgramData\'
    selection_immediate_execution:
        TaskContent|contains:
            - '<StartBoundary>'  # Task set to start immediately
        TaskContent|regex: '<StartBoundary>\d{4}-\d{2}-\d{2}T\d{2}:\d{2}:\d{2}'
    filter_system_tasks:
        SubjectUserName: 'SYSTEM'
        TaskName|startswith: '\Microsoft\'
    condition: >-
      selection_task_creation and (selection_high_privilege_indicators or selection_suspicious_actions or selection_suspicious_paths or selection_immediate_execution) and not filter_system_tasks
falsepositives:
    - Software installation creating maintenance task scheduling
    - System update mechanisms
level: medium

