# Linux/linux_systemd_user_service_persistence.yaml
title: Linux Systemd User Service Persistence Creation
id: b7e4a1c2-2001-4a5f-9c3d-bbbbbbbb0001
status: experimental
description: |
  Detects creation or enablement of user-level systemd units outside standard
  package management, a common persistence technique on Linux hosts.
author: zendannyy
date: 2026-06-06
references:
  - https://attack.mitre.org/techniques/T1543/002/
tags:
  - attack.persistence
  - attack.t1543.002
logsource:
  product: linux
  category: process_creation
detection:
  selection:
    Image|endswith:
      - '/systemctl'
      - '/systemd-run'
    CommandLine|contains:
      - '--user'
      - 'enable'
      - 'start'
  selection_write:
    Image|endswith:
      - '/cp'
      - '/mv'
      - '/tee'
      - '/install'
    CommandLine|contains:
      - '.config/systemd/user'
      - '/etc/systemd/system'
  condition: 1 of selection*
fields:
  - Image
  - CommandLine
  - User
falsepositives:
  - User-installed development services and legitimate desktop software
level: medium
